# btrborg.service
# vim: ft=systemd
[Unit]
Description=btrborg backup

[Service]
Type=oneshot
ExecStart=/usr/bin/btrborg-create
